# globus-connect/supervisord.conf
[supervisord]
nodaemon=true

[program:gridftp]
command=/usr/sbin/globus-gridftp-server -c /etc/gridftp.conf -pidfile /var/run/globus-gridftp-server.pid -no-detach -config-base-path /data


